# docker/ansible/motd/files/motd_handler.timer
[Unit]
Description=update motd daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
